# infra/examples/weekly-audit-report.service.example
[Unit]
Description=Verwaltung App weekly audit report
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__APP_USER__
Group=__APP_GROUP__
Environment=PHP_BIN=__PHP_BIN__
WorkingDirectory=__APP_ROOT__
ExecStart=/usr/bin/env bash __APP_ROOT__/infra/scripts/send-weekly-audit-report.sh --admin-email=__ADMIN_EMAIL__
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
